# supervisord/supervisord.conf
[supervisord]
nodaemon = true

[program:sshd]
command=/usr/sbin/sshd -D

[program:glances]
command=/root/env/bin/python3 -m glances -t 30 -C /root/conf/glances.conf --export influxdb2 -q